# docker/app/supervisord.conf
[program:php-fpm]
command=bash -c 'pkill -9 -f "^php-fpm: " 2>/dev/null; sleep 0.2; rm -f /run/php/php8.4-fpm.sock; exec /usr/sbin/php-fpm8.4 -F'
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
